# born2beroot/crontab
*/10 * * * * /home/lfrench/monitoring.sh
